# systemd/visa-esim-uz-bot.service
[Unit]
Description=Visa eSIM Uzbekistan Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/visa-esim-uz-bot
EnvironmentFile=/opt/visa-esim-uz-bot/.env
ExecStart=/usr/bin/python3 /opt/visa-esim-uz-bot/bot.py
Restart=always
RestartSec=5
User=visa-bot
Group=visa-bot

[Install]
WantedBy=multi-user.target

